# services/mongodb.service.conf
[program:mongodb]
user=docker
command=/usr/local/bin/run-mongo.sh
stdout_events_enabled = true
stderr_events_enabled = true
