# services/mongodb.service.conf
[program:mongodb]
user=root
command=/opt/mongodb/bin/mongod --httpinterface --rest --directoryperdb --dbpath /storage/mongodb/db
stdout_events_enabled = true
stderr_events_enabled = true
