# services/mongodb.service.conf
[program:mongodb]
user=root
command=/opt/mongodb/mongodb
stdout_events_enabled = true
stderr_events_enabled = true
